# deploy/systemd/opentunnel-relay.service
[Unit]
Description=OpenTunnel relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/opentunnel/relay.env
User=opentunnel
Group=opentunnel
ExecStart=/usr/local/bin/opentunnel relay --listen ${OPENTUNNEL_LISTEN} --public-url ${OPENTUNNEL_PUBLIC_URL} --artifact-path ${OPENTUNNEL_ARTIFACT_PATH} --version ${OPENTUNNEL_VERSION}
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
